# supervisor/conf.d/DEVELOPMENT/flask_app.conf
[program:flask_app]
command = python app.py
directory = /var/www/app
autostart = true
autorestart = true
startretries = 1
redirect_stderr = true
stdout_logfile = /tmp/flask_app.log
